# systemd/user/web-osint-candidate-model-downloads.service
[Unit]
Description=Web OSINT candidate model downloads
Documentation=file:/home/ops/dev/x-research/docs/DERIVED_ARCHITECTURE_IMPLEMENTATION_PLAN.md
Documentation=https://huggingface.co/docs/huggingface_hub/guides/cli
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ops/dev/x-research
Environment=WEB_OSINT_DATA_ROOT=/mnt/data/web-osint-platform
Environment=WEB_OSINT_MODELS_DIR=/mnt/data/web-osint-platform/models
Environment=WEB_OSINT_HF_VENV=/mnt/data/web-osint-platform/.venv-hf
Environment=WEB_OSINT_MODEL_LOG_DIR=/mnt/data/web-osint-platform/logs/candidate-model-downloads
Environment=WEB_OSINT_CANDIDATE_MODEL_MANIFEST=/mnt/data/web-osint-platform/model-download-manifests/candidate-models.tsv
Environment=WEB_OSINT_CANDIDATE_DOWNLOAD_JOBS=4
Environment=HF_HOME=/mnt/data/web-osint-platform/huggingface
Environment=HF_HUB_CACHE=/mnt/data/web-osint-platform/huggingface/hub
Environment=HF_XET_CACHE=/mnt/data/web-osint-platform/huggingface/xet
Environment=HF_TOKEN_FILE=/home/ops/dev/huggingface.md
Environment=HF_HUB_VERBOSITY=info
ExecStart=/home/ops/dev/x-research/scripts/download_candidate_models.sh
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
