# systemd/user/web-osint-rebrowser-launch-helper.service
[Unit]
Description=Web OSINT Rebrowser LaunchCapture helper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=WEB_OSINT_REPO_ROOT=/opt/web-osint-platform
Environment=WEB_OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=PANDAPROXY_URL=https://pandaproxy.example.invalid
Environment=REBROWSER_LAUNCH_HELPER_BIND=127.0.0.1
Environment=REBROWSER_LAUNCH_HELPER_PORT=18231
Environment=REBROWSER_DISPLAY=:96
Environment=REBROWSER_CDP_URL=http://127.0.0.1:9250
Environment=REBROWSER_CHROME_BIN=/path/to/chrome
Environment=REBROWSER_PROFILE_DIR=/srv/web-osint-platform/browser-profile
Environment=REBROWSER_SESSION_OPEN_URL=https://rebrowser-session.example.invalid/vnc.html
Environment=REBROWSER_LAUNCH_CAPTURE_MODE=rendered-web
Environment=REBROWSER_LAUNCH_SETTLE_MS=2500
Environment=REBROWSER_LAUNCH_NAV_TIMEOUT_MS=60000
Environment=PLAYWRIGHT_MODULE=/path/to/rebrowser-playwright-core/index.mjs
EnvironmentFile=-%h/.config/web-osint-platform/web-osint.env
ExecStart=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && exec python3 scripts/rebrowser_launch_helper.py'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
