# systemd/user/web-osint-embedding-worker.service
[Unit]
Description=Web OSINT embedding enrichment worker
Documentation=file:/home/ops/dev/x-research/docs/LOCAL_INFERENCE.md
After=network-online.target local-inference.service
Wants=network-online.target local-inference.service

[Service]
Type=simple
WorkingDirectory=/home/ops/dev/x-research/workers/embedding-worker
Environment=WEB_OSINT_DATA_ROOT=/mnt/data/web-osint-platform
Environment=WEB_OSINT_EMBEDDING_WORKER_VENV=/mnt/data/web-osint-platform/.venv-embedding-worker
Environment=WEB_OSINT_CPU_RESERVED_THREADS=2
Environment=REDPANDA_BROKERS=127.0.0.1:19092
Environment=REDPANDA_GROUP_ID=web-osint-embedding-worker-v1
Environment=KAFKA_BROKERS=127.0.0.1:19092
Environment=KAFKA_GROUP_ID=web-osint-embedding-worker-v1
Environment=LOCAL_INFERENCE_URL=http://127.0.0.1:18200
Environment=QDRANT_URL=http://127.0.0.1:16333
Environment=QDRANT_COLLECTION=web_osint_evidence_v1
Environment=HTTP_ADDR=127.0.0.1:18201
Environment=TOKENIZERS_PARALLELISM=false
ExecStartPre=/home/ops/dev/x-research/scripts/init_embedding_worker_venv.sh
ExecStart=/home/ops/dev/x-research/scripts/run_with_cpu_thread_guard.sh /mnt/data/web-osint-platform/.venv-embedding-worker/bin/python /home/ops/dev/x-research/workers/embedding-worker/embedding_worker.py
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
