# systemd/user/web-osint-embedding-worker.service
[Unit]
Description=Web OSINT embedding enrichment worker
Documentation=https://github.com/prmartinow/web-osint-platform/blob/main/docs/LOCAL_INFERENCE.md
After=network-online.target local-inference.service
Wants=network-online.target local-inference.service

[Service]
Type=simple
Environment=WEB_OSINT_REPO_ROOT=/opt/web-osint-platform
Environment=WEB_OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=WEB_OSINT_EMBEDDING_WORKER_VENV=/srv/web-osint-platform/.venv-embedding-worker
Environment=WEB_OSINT_CPU_RESERVED_THREADS=2
Environment=REDPANDA_BROKERS=kafka.example.invalid:9092
Environment=REDPANDA_GROUP_ID=web-osint-embedding-worker-v1
Environment=KAFKA_BROKERS=kafka.example.invalid:9092
Environment=KAFKA_GROUP_ID=web-osint-embedding-worker-v1
Environment=LOCAL_INFERENCE_URL=https://local-inference.example.invalid
Environment=QDRANT_URL=https://qdrant.example.invalid
Environment=QDRANT_COLLECTION=web_osint_evidence_v1
Environment=HTTP_ADDR=replace-with-bind-host:18201
Environment=TOKENIZERS_PARALLELISM=false
EnvironmentFile=-%h/.config/web-osint-platform/web-osint.env
ExecStartPre=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && scripts/init_embedding_worker_venv.sh'
ExecStart=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && exec scripts/run_with_cpu_thread_guard.sh "$WEB_OSINT_EMBEDDING_WORKER_VENV/bin/python" workers/embedding-worker/embedding_worker.py'
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
