# systemd/user/web-osint-media-ocr-worker.service
[Unit]
Description=Web OSINT media OCR client worker
Documentation=https://github.com/prmartinow/web-osint-platform/blob/main/docs/OPERATING.md
After=network-online.target local-inference.service web-osint-media-router.service
Wants=network-online.target local-inference.service

[Service]
Type=simple
Environment=WEB_OSINT_REPO_ROOT=/opt/web-osint-platform
Environment=WEB_OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=WEB_OSINT_MEDIA_ENRICHMENT_VENV=/srv/web-osint-platform/.venv-media-enrichment
Environment=OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=WEB_OSINT_CPU_RESERVED_THREADS=2
Environment=REDPANDA_BROKERS=kafka.example.invalid:9092
Environment=KAFKA_BROKERS=kafka.example.invalid:9092
Environment=CLICKHOUSE_URL=https://clickhouse.example.invalid
Environment=CLICKHOUSE_DATABASE=web_osint
Environment=CLICKHOUSE_USER=web_osint
Environment=LOCAL_INFERENCE_URL=https://local-inference.example.invalid
Environment=MEDIA_WORKER_ROLE=ocr
Environment=MEDIA_OCR_HTTP_ADDR=replace-with-bind-host:18212
Environment=TOKENIZERS_PARALLELISM=false
EnvironmentFile=-%h/.config/web-osint-platform/web-osint.env
ExecStartPre=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && scripts/init_media_enrichment_venv.sh'
ExecStart=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && exec scripts/run_with_cpu_thread_guard.sh "$WEB_OSINT_MEDIA_ENRICHMENT_VENV/bin/python" workers/media-enrichment/media_enrichment_worker.py ocr'
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
