# systemd/user/web-osint-media-ocr-worker.service
[Unit]
Description=Web OSINT media OCR client worker
Documentation=file:/home/ops/dev/x-research/docs/OPERATING.md
After=network-online.target local-inference.service web-osint-media-router.service
Wants=network-online.target local-inference.service

[Service]
Type=simple
WorkingDirectory=/home/ops/dev/x-research
EnvironmentFile=-/home/ops/dev/x-research/.env
Environment=OSINT_DATA_ROOT=/mnt/data/x-research
Environment=WEB_OSINT_CPU_RESERVED_THREADS=2
Environment=REDPANDA_BROKERS=127.0.0.1:19092
Environment=KAFKA_BROKERS=127.0.0.1:19092
Environment=CLICKHOUSE_URL=http://127.0.0.1:18123
Environment=CLICKHOUSE_DATABASE=web_osint
Environment=CLICKHOUSE_USER=web_osint
Environment=LOCAL_INFERENCE_URL=http://127.0.0.1:18200
Environment=MEDIA_WORKER_ROLE=ocr
Environment=MEDIA_OCR_HTTP_ADDR=127.0.0.1:18212
Environment=TOKENIZERS_PARALLELISM=false
ExecStartPre=/home/ops/dev/x-research/scripts/init_media_enrichment_venv.sh
ExecStart=/home/ops/dev/x-research/scripts/run_with_cpu_thread_guard.sh /mnt/data/web-osint-platform/.venv-media-enrichment/bin/python /home/ops/dev/x-research/workers/media-enrichment/media_enrichment_worker.py ocr
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
