# systemd/user/web-osint-webpage-extraction-worker.service
[Unit]
Description=Web OSINT webpage extraction worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=WEB_OSINT_REPO_ROOT=/opt/web-osint-platform
Environment=WEB_OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=WEB_OSINT_WEBPAGE_EXTRACTION_VENV=/srv/web-osint-platform/.venv-webpage-extraction
Environment=OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=KAFKA_BROKERS=kafka.example.invalid:9092
Environment=PANDAPROXY_URL=https://pandaproxy.example.invalid
Environment=WEBPAGE_EXTRACTION_HTTP_ADDR=replace-with-bind-host:18221
Environment=WEBPAGE_EXTRACTION_REQUEST_TIMEOUT=45
Environment=WEBPAGE_EXTRACTION_MAX_HTML_BYTES=8388608
EnvironmentFile=-%h/.config/web-osint-platform/web-osint.env
ExecStartPre=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && scripts/init_webpage_extraction_venv.sh'
ExecStart=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && exec "${WEB_OSINT_WEBPAGE_EXTRACTION_PYTHON:-${WEB_OSINT_WEBPAGE_EXTRACTION_VENV:?set WEB_OSINT_WEBPAGE_EXTRACTION_VENV}/bin/python}" workers/webpage-extraction/webpage_extraction_worker.py run'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
