# systemd/user/web-osint-qdrant-embedding-backfill.service
[Unit]
Description=Web OSINT one-shot Qdrant embedding backfill
Documentation=https://github.com/prmartinow/web-osint-platform/blob/main/docs/LOCAL_INFERENCE.md
After=network-online.target local-inference.service
Wants=network-online.target local-inference.service

[Service]
Type=simple
EnvironmentFile=-%h/.config/web-osint-platform/web-osint.env
Environment=WEB_OSINT_REPO_ROOT=/opt/web-osint-platform
Environment=WEB_OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=WEB_OSINT_EMBEDDING_WORKER_VENV=/srv/web-osint-platform/.venv-embedding-worker
Environment=WEB_OSINT_CPU_RESERVED_THREADS=2
Environment=CLICKHOUSE_URL=https://clickhouse.example.invalid
Environment=CLICKHOUSE_DATABASE=web_osint
Environment=CLICKHOUSE_USER=web_osint
Environment=LOCAL_INFERENCE_URL=https://local-inference.example.invalid
Environment=QDRANT_URL=https://qdrant.example.invalid
Environment=QDRANT_COLLECTION=web_osint_evidence_v1
Environment=BACKFILL_BATCH_SIZE=1
ExecStartPre=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && scripts/init_embedding_worker_venv.sh'
ExecStart=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && exec scripts/run_with_cpu_thread_guard.sh "$WEB_OSINT_EMBEDDING_WORKER_VENV/bin/python" scripts/backfill_qdrant_embeddings.py'
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
