# systemd/user/web-osint-media-router.service
[Unit]
Description=Web OSINT media enrichment router
Documentation=https://github.com/prmartinow/web-osint-platform/blob/main/docs/OPERATING.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=WEB_OSINT_REPO_ROOT=/opt/web-osint-platform
Environment=WEB_OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=WEB_OSINT_MEDIA_ENRICHMENT_VENV=/srv/web-osint-platform/.venv-media-enrichment
Environment=OSINT_DATA_ROOT=/srv/web-osint-platform/data
Environment=REDPANDA_BROKERS=kafka.example.invalid:9092
Environment=KAFKA_BROKERS=kafka.example.invalid:9092
Environment=CLICKHOUSE_URL=https://clickhouse.example.invalid
Environment=CLICKHOUSE_DATABASE=web_osint
Environment=CLICKHOUSE_USER=web_osint
Environment=MEDIA_WORKER_ROLE=router
Environment=MEDIA_ROUTER_HTTP_ADDR=replace-with-bind-host:18211
Environment=TOKENIZERS_PARALLELISM=false
EnvironmentFile=-%h/.config/web-osint-platform/web-osint.env
ExecStartPre=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && scripts/init_media_enrichment_venv.sh'
ExecStart=/usr/bin/env bash -lc 'cd "${WEB_OSINT_REPO_ROOT:?set WEB_OSINT_REPO_ROOT}" && exec "$WEB_OSINT_MEDIA_ENRICHMENT_VENV/bin/python" workers/media-enrichment/media_enrichment_worker.py router'
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
